# sync-shelf-tags.timer
# =============================================================================
#  sync-shelf-tags.timer  — runs the sync periodically.
#  Default: every 15 minutes, plus 2 min after boot. Adjust OnUnitActiveSec
#  to taste (e.g. 1h if you don't need it snappy).
# =============================================================================

[Unit]
Description=Run shelf->tag sync periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
# don't pile up runs if one is slow; and run a missed one after downtime
Persistent=true

[Install]
WantedBy=timers.target
