# scripts/sweep-books.service
[Unit]
Description=gunit book cycle: sweep-import missed files, drain tag queue, sync tags to shelves
After=docker.service watch-downloads.service
Wants=docker.service

[Service]
Type=oneshot
# runs as root, like watch-downloads.service, so it can docker exec + move files
User=root
# 1) import any files the live watcher missed, then 2) tag whatever's now in the
# library that's waiting in the tag queue, then 3) mirror configured tags onto
# their calibre-web shelves. Sequential and in dependency order: a book can go
# import -> tag -> shelf in a single cycle. The '-' on steps 2 and 3 means an
# earlier step's failure doesn't abort the later ones.
ExecStart=/home/robmorgan/gunit/sweep-books.sh
ExecStart=-/home/robmorgan/gunit/tag-queue.sh
ExecStart=-/home/robmorgan/gunit/sync-tag-to-shelf.sh
# don't let a slow sweep run forever; the single-instance flock also guards overlap
TimeoutStartSec=600
# log to the journal (view with: journalctl -u sweep-books.service)
StandardOutput=journal
StandardError=journal
